# files/insights-client.timer.d/override.conf
[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=14400
OnBootSec=0min
